# GFWTools/V2Ray/upstart/v2ray.conf
description "V2Ray"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

respawn

script
	cd /root/v2ray-v3.50-linux-64 && ./v2ray
end script
