# Docker/docker.service.d/proxy.conf
[Service]
Environment="ALL_PROXY=socks5://192.168.2.1:1080"
